# rollback-on-backup.service
[Unit]
Description=Rollback on backup root to restore live root, then switch back
DefaultDependencies=yes
After=local-fs.target
ConditionPathExists=/boot/firmware/rollback.flag

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/rollback-on-backup.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
